# rootfs/etc/s6-overlay/s6-rc.d/tautulli/run
#!/usr/bin/with-contenv sh

CONFIG_FILE=/config/config.ini

# create initial config file
if [[ ! -f "$CONFIG_FILE" ]]; then
  cat <<EOF >"$CONFIG_FILE"
[General]
check_github = 0
[Advanced]
system_analytics = 0
EOF
  chown $PUID:$PGID "$CONFIG_FILE"
else
  # disable check_github for updates
  sed -i "/^check_github[ ]*/s/=.*/= 0/" "$CONFIG_FILE"
  # disable analytics reports
  sed -i "/^system_analytics[ ]*/s/=.*/= 0/" "$CONFIG_FILE"
fi

# change chdir
cd /config

# run with limited permissions
exec s6-setuidgid $PUID:$PGID python3 /app/Tautulli.py --nofork --nolaunch --datadir=/config
